# docker/entrypoint.conf
[program:entrypoint]
command=/usr/local/bin/entrypoint.sh
directory=/var/www/html/
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0
